# server/nginx-s6/1.19.1/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh

exec 2>&1

exec nginx -c /etc/nginx/nginx.conf
